# crontab
#m h dom mon dow command
*/10 * * * * sh /home/abhishek/scripts/battery-manager.sh >/dev/null 2>&1